# deploy/systemd/mcharness-cockpit.service
[Unit]
Description=McHarness public cockpit backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/mcharness-public-export
Environment=PYTHONPATH=/root/mcharness-public-export
Environment=PYTHONUNBUFFERED=1
Environment=MCHARNESS_DATA_ROOT=/var/lib/mcharness-cockpit
Environment=MCHARNESS_PUBLIC_WRITE_ENABLED=false
ExecStart=/root/hybrid-agent-os/.venv/bin/python -m uvicorn src.server.api:app --host 127.0.0.1 --port 8124 --log-level warning
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
